# x86_64/var/gentoo/ofs-boot/user/etc/local.d/50-x.start
[ "$(cat /proc/cmdline | grep nox)" == "" ] && /usr/local/bin/x.start
