# ofs-boot/base/etc/init.d/empty-vt
#!/sbin/openrc-run

depend()
{
	after mount-ro
}

start()
{
	chvt 15
}
