# ofs-boot/base/etc/profile.d/pycache.sh
export PYTHONPYCACHEPREFIX=/var/tmp/pycache
