# arm64/var/gentoo/ofs-boot-build-kodi/misc/etc/profile.d/emscripten.sh
export EMSDK=/usr/lib64/emscripten
export TMPDIR=/var/tmp
