# ofs-boot/misc/etc/local.d/15-alsaloop.start
alsaloop -C "hw:Loopback,1" -P "plug:equal" -t 15000 -d
